# StopNvidiaInArch.service
[Unit]
Description=Nvidia Script

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/StopNvidiaInArch.sh

[Install]
WantedBy=multi-user.target
